# win_agent_createtask.xml
<?xml version="1.0" ?>
<!--
This schedules the CABS agent as a task on windows, so it will run regardless of if a user is logged on or not.
It will start when the system is booted.
-->
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
	<RegistrationInfo>
		<Date>2011-05-19</Date>
		<Author>JordanWhitehead</Author>
		<Version>0.5.0</Version>
		<Description>Starts CABS Agent on system boot.</Description>
	</RegistrationInfo>
	<Triggers>
		<BootTrigger>
			<Enabled>true</Enabled>
		</BootTrigger>
	</Triggers>
	<Principals>
		<Principal>
			<UserId>Administrator</UserId>
			<LogonType>5</LogonType>
		</Principal>
	</Principals>
	<Settings>
		<Enabled>true</Enabled>
		<AllowStartOnDemand>true</AllowStartOnDemand>
		<AllowHardTerminate>true</AllowHardTerminate>
	</Settings>
	<Actions>
		<Exec>
			<Command></Command>
		</Exec>
	</Actions>
</Task>
